# shared/config/systemd/openclaw-bot@.service
[Unit]
Description=Bot Fleet (OpenClaw) - %i
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/Bot-Fleet-Inc/fleet-ops
StartLimitBurst=5
StartLimitIntervalSec=600

[Service]
Type=simple
User=bot
Group=bot
WorkingDirectory=/opt/bot/workspace/%i
EnvironmentFile=/opt/bot/secrets/%i.env
Environment=HOME=/opt/bot
Environment=OPENCLAW_HOME=/opt/bot/workspace/%i/.openclaw
ExecStart=/usr/local/bin/openclaw gateway run
Restart=always
RestartSec=10

# Faster crash recovery (from openclaw-config production hardening)
# Default systemd throttle is 10s; 5s reduces downtime on repeated restarts
ThrottleInterval=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/bot /var/log/bot
ProtectHome=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bot-%i

# Environment expectations (documented, not enforced here):
#   GITHUB_TOKEN, ANTHROPIC_API_KEY, GEMINI_API_KEY,
#   OPENCLAW_HOOK_TOKEN, TELEGRAM_BOT_TOKEN,
#   CHAT_WORKER_TOKEN (legacy), BOT_NAME, LOCAL_LLM_URL

[Install]
WantedBy=multi-user.target
